# content/service/wst/run
#!/bin/sh

exec 2>&1
exec wst server ws://[::]:3001
