# content/service/wst/run
#!/bin/sh

exec 2>&1
exec appwst server ws://[::]:3001
